# bin/hommypay_trx/app/hmpy-app.service
[Unit]
Description=Hommypay App Service

[Service]
Type=simple
Restart=always
RestartSec=5s
WorkingDirectory=/home/rm/projects/hommypay_trx/app
ExecStart=/home/rm/projects/hommypay_trx/app/hommypay-trx-app

[Install]
WantedBy=multi-user.target