# bin/hommypay_trx/sale/hmpy-sale.service
[Unit]
Description=Hommypay Sale Service

[Service]
Type=simple
Restart=always
RestartSec=5s
WorkingDirectory=/home/rm/projects/hommypay_trx/sale
ExecStart=/home/rm/projects/hommypay_trx/sale/hommypay-trx-sale

[Install]
WantedBy=multi-user.target